# config/supervisor/vnc.conf
[program:vnc]
user=root
command=x11vnc -forever -rfbport 5900 -noxrecord -noxfixes -noxdamage -rfbauth /etc/vncpwd -display :99
autostart=true
autorestart=true